# root/etc/systemd/system/asterisk.service.d/cdr.conf
# Start Asterisk after MySQL to avoid CDR failure
[Unit]
After=mysqld.service
Wants=mysqld.service

